# scripts/backup/nutriflow-backup.timer
[Unit]
Description=NutriFlow PostgreSQL backup — daily at 03:00 UTC

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
Unit=nutriflow-backup.service

[Install]
WantedBy=timers.target
